# deploy/invisiblego.service
[Unit]
Description=InvisibleGo web server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=invisiblego
Group=invisiblego
WorkingDirectory=/opt/invisiblego
Environment="PYTHONUNBUFFERED=1"
ExecStart=/opt/invisiblego/venv/bin/uvicorn transport.web.server:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
